# deploy/systemd/status-saver-rotate.service
[Unit]
Description=status-saver retention rotation
Documentation=https://github.com/ppoloczek/status-saver

[Service]
Type=oneshot
User=status-saver
Group=status-saver
ExecStart=/usr/local/bin/status-saver rotate --config /etc/status-saver/config.yaml

StateDirectory=status-saver
StateDirectoryMode=0750

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/var/lib/status-saver
