# deploy/launchd/com.github.peetzweg.status-saver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS LaunchAgent for status-saver.

  Install: copy to ~/Library/LaunchAgents/ and substitute YOURNAME below.
  Load:    launchctl load ~/Library/LaunchAgents/com.github.peetzweg.status-saver.plist

  Untested end-to-end. If you get it working (or it breaks), please report
  back at https://github.com/peetzweg/status-saver/issues.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.github.peetzweg.status-saver</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOURNAME/.local/bin/status-saver</string>
        <string>run</string>
        <string>--config</string>
        <string>/Users/YOURNAME/Library/Application Support/status-saver/config.yaml</string>
    </array>

    <!-- Start at login and keep running. launchd will restart the process
         if it exits non-zero (KeepAlive with SuccessfulExit=false). If
         status-saver exits 1 because of a permanent logout, the restart
         loop will keep hammering — watch the log and re-pair manually. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- macOS may idle the process when the Mac is locked / on battery.
         LegacyTimers reduces (doesn't eliminate) that. For a truly always-on
         service, run on a Linux server instead or switch to a LaunchDaemon. -->
    <key>LegacyTimers</key>
    <true/>

    <!-- Logs go here. tail -f to follow. -->
    <key>StandardOutPath</key>
    <string>/Users/YOURNAME/Library/Logs/status-saver.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOURNAME/Library/Logs/status-saver.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/YOURNAME/Library/Application Support/status-saver</string>

    <!-- Process stays nice; keep CPU priority normal. Tune if the box
         is otherwise idle. -->
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
